# addons/copilot_bridge/rootfs/etc/services.d/copilot_bridge/run
#!/usr/bin/with-contenv bashio

export BRIDGE_PORT=8099
export BRIDGE_API_KEY="$(bashio::config 'bridge_api_key')"
export GITHUB_TOKEN="$(bashio::config 'github_token')"
export GITHUB_OAUTH_CLIENT_ID="$(bashio::config 'github_oauth_client_id')"
export GITHUB_OAUTH_SCOPES="$(bashio::config 'github_oauth_scopes')"
export ASSISTANT_PROFILE="$(bashio::config 'assistant_profile')"
export READ_ONLY_MODE="$(bashio::config 'read_only_mode')"
export ALLOW_HOME_ASSISTANT_ACTIONS="$(bashio::config 'allow_home_assistant_actions')"
export ALLOW_FILESYSTEM_ACCESS="$(bashio::config 'allow_filesystem_access')"
export ENABLE_INTEGRATION_DISCOVERY="$(bashio::config 'enable_integration_discovery')"
export ENABLE_HACS_DISCOVERY="$(bashio::config 'enable_hacs_discovery')"
export ENABLE_TOOLING_DISCOVERY="$(bashio::config 'enable_tooling_discovery')"
export ENABLE_ZEROCONF_DISCOVERY="$(bashio::config 'enable_zeroconf_discovery')"
export ENABLE_HOME_ASSISTANT_MCP="$(bashio::config 'enable_home_assistant_mcp')"
export HOME_ASSISTANT_MCP_URL="$(bashio::config 'home_assistant_mcp_url')"
export HOME_ASSISTANT_MCP_BEARER_TOKEN="$(bashio::config 'home_assistant_mcp_bearer_token')"
export HOME_ASSISTANT_MCP_API_KEY="$(bashio::config 'home_assistant_mcp_api_key')"
export ALLOWED_PATHS="$(bashio::config 'allowed_paths')"
export LOG_LEVEL="$(bashio::config 'log_level')"

exec python3 /app/server.py
